# rpm/networking-vpp-agent.service
[Unit]
Description=Networking VPP ML2 Agent
Requires=vpp.service
After=vpp.service

[Service]
ExecStart=/usr/bin/vpp-agent --config-file /etc/neutron/plugins/ml2/vpp_agent.ini
Type=simple

[Install]
WantedBy=multi-user.target

